# mockserver.service
[Unit]
Description=Tax App Mock Server (Production)
After=network.target

[Service]
User=taxapp
Group=taxapp
WorkingDirectory=/opt/taxapp
ExecStart=/usr/bin/python3 /opt/taxapp/mockserver_production.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

# Log to syslog
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=mockserver

[Install]
WantedBy=multi-user.target 