# mockserver.service
[Unit]
Description=Tax App Mock Server (Production)
After=network.target

[Service]
User=taxapp
Group=taxapp
WorkingDirectory=/opt/taxapp
# Modified to run without SSL since Apache handles that
ExecStart=/usr/bin/python3 -c "import sys; sys.argv = ['', '--no-ssl']; exec(open('/opt/taxapp/mockserver_production.py').read())"
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

# Log to syslog
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=mockserver

[Install]
WantedBy=multi-user.target 